# Setup/com.christian.secondbrain.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.christian.secondbrain.telegram-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/christian/.second-brain/telegram-bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/christian/Documents/Codex/Second brain</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SECOND_BRAIN_VAULT</key>
    <string>/Users/christian/Documents/Codex/Second brain</string>
    <key>SECOND_BRAIN_COMMAND_SCRIPT</key>
    <string>/Users/christian/.second-brain/second-brain-command.sh</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/second-brain-telegram-bot.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/second-brain-telegram-bot.err</string>
</dict>
</plist>
